# launchd/com.safetybid.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.safetybid.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "/Users/dev06/Safety diagnosis bid" &amp;&amp; ./.venv/bin/python -m src.monitor</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dev06/Safety diagnosis bid</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/dev06/Safety diagnosis bid/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dev06/Safety diagnosis bid/logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>ko_KR.UTF-8</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
